# rosec/src/enable/templates/rosecd.service
[Unit]
Description=rosec - multi-provider Secret Service daemon
Documentation=https://github.com/jmylchreest/rosec
After=dbus.service
Requires=dbus.service
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=org.freedesktop.secrets
ExecStart=@@ROSECD@@
Restart=on-failure
RestartSec=5

# Security hardening.
#
# The SSH agent exposes a FUSE filesystem at $XDG_RUNTIME_DIR/rosec/ssh/ that
# must be visible in the user's session mount namespace.  Any systemd directive
# that implies a private mount namespace (ProtectSystem, PrivateTmp,
# ProtectKernelTunables, ProtectKernelModules, ProtectControlGroups) would make
# the FUSE mount invisible outside the service, so they are intentionally
# omitted.
#
# NoNewPrivileges is omitted because fusermount3's setuid bit is required.
# MemoryDenyWriteExecute is omitted because wasmtime's JIT needs W+X pages.
# RestrictNamespaces is omitted because fusermount3 may use mount namespaces.
DeviceAllow=/dev/fuse rw
# Allow mlockall so decrypted secrets are never swapped to disk.
LimitMEMLOCK=infinity
RestrictRealtime=yes
PrivateNetwork=no
Environment=RUST_LOG=info

[Install]
WantedBy=graphical-session.target
Also=rosecd.socket
